# .config/systemd/user/lei-work.service
[Unit]
Description=lei synchronize

[Service]
Type=oneshot
ExecStart=lei up --all
